# nuc-full/meshcore-export.timer
[Unit]
Description=Run the MeshCore triangulator export every hour

[Timer]
# First run 10 min after boot (give the collector time to connect), then hourly.
OnBootSec=10min
OnUnitActiveSec=1h
# If the NUC was off when a run was due, run once on the next boot.
Persistent=true

[Install]
WantedBy=timers.target
