# nuc-full/meshcore-triangulator.service
[Unit]
Description=MeshCore MQTT Triangulator collector
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=mctri
Group=mctri
WorkingDirectory=/opt/meshcore-triangulator
ExecStart=/opt/meshcore-triangulator/.venv/bin/python /opt/meshcore-triangulator/collector.py --config /opt/meshcore-triangulator/config.ini
Restart=always
RestartSec=10
# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/meshcore-triangulator

[Install]
WantedBy=multi-user.target
