# nuc-full/meshcore-validate.service
[Unit]
Description=MeshCore triangulator — accuracy self-check (locate.py --validate)
After=meshcore-triangulator.service
Wants=meshcore-triangulator.service

[Service]
Type=oneshot
User=mctri
Group=mctri
WorkingDirectory=/opt/meshcore-triangulator
# Validate against every GPS-known node and write a fresh, timestamped report.
# Latest report:  /opt/meshcore-triangulator/validate-report.txt
# Trend log:      /opt/meshcore-triangulator/validate-history.log  (the aggregate block, appended)
ExecStart=/bin/sh -c '\
  { echo "=== $(date -u +%%Y-%%m-%%dT%%H:%%M:%%SZ) ==="; \
    /opt/meshcore-triangulator/.venv/bin/python /opt/meshcore-triangulator/locate.py \
      --config /opt/meshcore-triangulator/config.ini --validate; \
  } > /opt/meshcore-triangulator/validate-report.txt.tmp 2>&1 && \
  mv /opt/meshcore-triangulator/validate-report.txt.tmp /opt/meshcore-triangulator/validate-report.txt; \
  { echo "=== $(date -u +%%Y-%%m-%%dT%%H:%%M:%%SZ) ==="; \
    grep -A6 "Aggregate over" /opt/meshcore-triangulator/validate-report.txt; \
    echo; } >> /opt/meshcore-triangulator/validate-history.log'
TimeoutStartSec=3600
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/meshcore-triangulator

[Install]
WantedBy=multi-user.target
