# nuc-full/meshcore-export.service
[Unit]
Description=MeshCore triangulator — locate all nodes and export JSON
After=meshcore-triangulator.service
Wants=meshcore-triangulator.service

[Service]
Type=oneshot
User=mctri
Group=mctri
WorkingDirectory=/opt/meshcore-triangulator
# Locate every candidate node and write the up-to-date location JSON.
# Output includes per node: name, last_seen (last reception, UTC) and last_advert.
ExecStart=/opt/meshcore-triangulator/.venv/bin/python /opt/meshcore-triangulator/export_triangulator.py \
  --repo /opt/meshcore-triangulator \
  --out /opt/meshcore-triangulator/web/triangulator-targets.json \
  --tier-max 1
# A full network pass can take 15-25 min; don't let systemd kill it early.
TimeoutStartSec=3600
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/meshcore-triangulator
